# scripts/launchd/pl.octadecimal.knowledge-sync-m1.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>pl.octadecimal.knowledge-sync-m1</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@REPO_ROOT@/scripts/octa-knowledge-sync-m1.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>@REPO_ROOT@</string>
  <key>RunAtLoad</key>
  <true/>
  <key>WatchPaths</key>
  <array>
    <string>@KNOWLEDGE_ROOT@</string>
  </array>
  <key>ThrottleInterval</key>
  <integer>120</integer>
  <key>StandardOutPath</key>
  <string>@HOME@/.octa/logs/knowledge-sync-m1-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/.octa/logs/knowledge-sync-m1-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>KNOWLEDGE_ROOT</key>
    <string>@KNOWLEDGE_ROOT@</string>
    <key>M1_SSH_HOST</key>
    <string>m1-ceo</string>
  </dict>
</dict>
</plist>
