# scripts/launchd/pl.octadecimal.embed-knowledge-dev.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>pl.octadecimal.embed-knowledge-dev</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@REPO_ROOT@/scripts/octa-knowledge-sync-dev.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>@REPO_ROOT@</string>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>@HOME@/.octa/logs/embed-sync-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/.octa/logs/embed-sync-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>KNOWLEDGE_ROOT</key>
    <string>@HOME@/Developer/Knowledge</string>
  </dict>
</dict>
</plist>
